# deploy/server/systemd/ai-data-platform-web.service
[Unit]
Description=AI Data Platform Web
After=network-online.target ai-data-platform-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/ai-data-platform/apps/web
EnvironmentFile=/etc/ai-data-platform/ai-data-platform.env
ExecStart=/bin/sh -lc 'PORT=${WEB_PORT:-3000} HOST=${WEB_HOST:-127.0.0.1} HOSTNAME=${WEB_HOST:-127.0.0.1} /usr/local/bin/node server.js'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
