# deploy/server/systemd/ai-data-platform-api.service
[Unit]
Description=AI Data Platform API
After=network-online.target ai-data-platform-model-bridge.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/ai-data-platform/apps/api
EnvironmentFile=/etc/ai-data-platform/ai-data-platform.env
ExecStart=/bin/sh -lc 'PORT=${API_PORT:-3100} HOST=${API_HOST:-127.0.0.1} /usr/local/bin/node dist/server.js'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
