# deploy/server/systemd/ai-data-platform-control-plane-api.service
# Legacy migration reference only.
# Do not install this unit for new environments. Shared platform services now live in `home`.
[Unit]
Description=LEGACY AI Data Platform Control Plane API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/ai-data-platform/apps/control-plane-api
EnvironmentFile=/srv/ai-data-platform/deploy/server/ai-data-platform.env
ExecStart=/bin/sh -lc 'PORT=${CONTROL_PLANE_API_PORT:-3210} HOST=${CONTROL_PLANE_API_HOST:-0.0.0.0} /usr/local/bin/node dist/server.js'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
